# .config/systemd/user/dotlink.timer
[Unit]
Description=Run dotlink script every hour

[Timer]
OnUnitActiveSec=1h
Unit=dotlink.service

[Install]
WantedBy=timers.target
