# .config/systemd/user/dotlink.service
[Unit]
Description=Run dotlink script

[Service]
Type=oneshot
ExecStart=/home/y0usaf/dotfiles/scripts/dotlink.sh

[Install]
WantedBy=default.target
