# .config/systemd/user/paru-update.service
[Unit]
Description=Update Arch Linux packages with paru

[Service]
Type=oneshot
ExecStart=/usr/bin/paru

[Install]
WantedBy=default.target
