# .config/systemd/user/remove-cache.service
[Unit]
Description=Remove ~/.cache directory

[Service]
Type=oneshot
ExecStart=/bin/bash -c "rm -rf ~/.cache"

[Install]
WantedBy=multi-user.target
